# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:app]
command=python /app/main.py

[program:streamlit]
directory=/app/examples/ui/
command=streamlit run app.py
